# etc/supervisord.d/ceph-mds.ini
[program:ceph-mds]
command=/opt/qnib/ceph/mds/bin/start.sh
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart = true
