# etc/supervisord.d/ceph-update.ini
[program:ceph-update]
command=consul-template -consul consul.service.consul:8500 -template "/etc/consul-templates/templates/ceph.conf.ctmpl:/etc/ceph/ceph.conf:supervisorctl restart ceph-mon"
redirect_stderr=true
stdout_logfile=syslog
